# Go/src/HouseholdBot/HouseholdBot.conf
[program:HouseholdBot]
command=/home/orangepi/go/HouseholdBot
aurostart=true
autorestart=true
stderr_logfile=/var/log/telegramBot.err.log
stdout_logfile=/var/log/telegramBot.out.log
